# provision/vault.service
[Unit]
Description=Hashicorp Vault

[Service]
ExecStart=vault server -config=/vagrant/provision/vault-config.hcl
